# image_files/etc/supervisor/spider.conf
;route 
[program:cdspider_route] 
command=/usr/local/bin/cdspider_admin -c config/main.json route
process_name=%(program_name)s_%(process_num)02d 
numprocs=1 
directory=/app/loeye
umask=022 
priority=999 
autostart=true 
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

;schedule_rpc 
[program:cdspider_schedule_rpc] 
command=/usr/local/bin/cdspider_admin -c config/main.json schedule-rpc
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app/loeye
umask=022
priority=999
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

;schedule
[program:cdspider_schedule]
command=/usr/local/bin/cdspider_admin -c config/main.json schedule
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app/loeye
umask=022
priority=999
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

;fetch
[program:cdspider_fetch]
command=/usr/local/bin/cdspider_admin -c config/main.json fetch
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app/loeye
umask=022
priority=999
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

;spider_rpc
[program:cdspider_spider_rpc]
command=/usr/local/bin/cdspider_admin -c config/main.json spider-rpc
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app/loeye
umask=022
priority=999
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

;work
[program:cdspider_newtask_work]
command=/usr/local/bin/cdspider_admin -c config/main.json newtask-work
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app/loeye
umask=022
priority=999
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:cdspider_result_work]
command=/usr/local/bin/cdspider_admin -c config/main.json result-work
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app/loeye
umask=022
priority=999
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:cdspider_status_work]
command=/usr/local/bin/cdspider_admin -c config/main.json status-work
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app/loeye
umask=022
priority=999
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

;tool
[program:cdspider_work_tool]
command=/usr/local/bin/cdspider_admin -c config/main.json tool --tool-cls cdspider.tools.test_tool.test_tool
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app/loeye
umask=022
priority=999
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

;webui
[program:cdspider_web]
command=/usr/local/bin/cdspider_admin -c config/main.json web --scheduler-rpc http://sp-srpc.dev.loeyae.com --spider-rpc http://sp-rpc.dev.loeyae.com --need-auth --auth-type header --username admin --password 2021@Bys
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/app/loeye
umask=022
priority=999
autostart=true
autorestart=true
stdout_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/tmp/supervisor/%(program_name)s_%(process_num)02d_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10